# templates/default/onlinefs.service.erb
[Unit]
Description = Online feature store service.
After = network-online.target docker.service <%= @local_dependencies %>
Requires = docker.service <%= @local_dependencies %>

[Service]
<% if node['onlinefs']['consumer_config_dir'].nil? && node['onlinefs']['producer_config_dir'].nil? -%>
ExecStartPre=/bin/bash -c 'source <%= node['consul']['bin_dir']%>/domain_utils.sh && _is_service_ready "kafka" "<%= @kafka_fqdn %>"'
<% end -%>
ExecStart = /usr/bin/docker run --rm --name %n \
           -v <%= node['onlinefs']['home'] %>:/onlinefs \
           -v <%= node['onlinefs']['data_volume']['root_dir']%>:<%= node['onlinefs']['data_volume']['root_dir']%> \
           -v <%= @crypto_dir %>:/certs \
           -p <%= node['onlinefs']['monitoring'] %>:<%= node['onlinefs']['monitoring'] %> \
           --network=host \
           -e ONLINEFS_HOME=/onlinefs \
           -e JAVA_OPTS="-Xms<%= node['onlinefs']['java_start_heap_size'] %> -Xmx<%= node['onlinefs']['java_max_heap_size'] %>" \
           <%= @image_name %>
ExecStop = /usr/bin/docker stop %n 

Restart = on-failure
RestartSec = 2

[Install]
WantedBy = multi-user.target
